# autodrive-canlog.service
[Unit]
Description=AutoDrive hourly CAN logger
After=default.target

[Service]
Type=simple
ExecStartPre=/usr/bin/mkdir -p /home/oxbo/data/can
ExecStart=%h/c/autodrive/hourly_can_logger.py /home/oxbo/data/can --can-bus can0 --flush-every 1
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
