# autodrive-api.service
[Unit]
Description=AutoDrive REST API server
After=default.target

[Service]
Type=simple
ExecStart=%h/c/autodrive/api_server.py --host 0.0.0.0 --port 8080 --can-bus can0
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
